# buildtools/bin/crontab
30 2 * * 1 ~/Sety-project/devops/buildtools/bin/pyrun_static.sh
5 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_all.sh ftx debug
30 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_restart_tradeexecutor.sh ftx debug
