# buildtools/bin/crontab
30 2 * * 1 pyrun_static.sh
5 * * * * pyrun_all.sh
