# buildtools/bin/crontab
30 2 * * 1 ~/Sety-project/devops/buildtools/bin/pyrun_static.sh
5 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_histfeed.sh
10-60/15 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_riskpnl.sh ftx debug
10-60/15 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_riskpnl.sh ftx SysPerp
10 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_pfoptimizer.sh ftx debug
10 * * * * ~/Sety-project/devops/buildtools/bin/pyrun_pfoptimizer.sh ftx SysPerp
